# docker/web-sa/conf/supervisor/nginx.ini
[program:nginx]
command=/usr/sbin/nginx
autorestart=true
logfile_maxbytes=50MB
logfile_backup=2
stdout_logfile=/var/log/nginx/%(program_name)s.log
stderr_logfile=/var/log/nginx/%(program_name)s.log
